# scripts/doh-daemon@.service
[Unit]
Description=DOH Auto-commit Monitoring Daemon
Documentation=https://github.com/your-username/doh
After=network.target

[Service]
Type=oneshot
User=%i
Group=%i
ExecStart=/usr/local/bin/doh daemon --once
StandardOutput=journal
StandardError=journal
SyslogIdentifier=doh-daemon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/.doh

[Install]
WantedBy=multi-user.target
